# services/files/updaterepos.timer
[Timer]
OnCalendar=*-*-* *:00,20,40:00

[Install]
WantedBy=timers.target
