# services/files/updaterepos.service
[Service]
ExecStart=/bin/bash ${HOME}/.local/bin/updaterepos
Type=oneshot

[Install]
WantedBy=default.target
